# services/battery-monitor.service
[Unit]
Description=Battery Monitor — MATE Notifications
# On attend simplement que la session utilisateur de base soit prête
After=default.target

[Service]
Type=simple
# Utilise le chemin exact que vous avez testé avec succès
ExecStart=%h/.local/bin/battery_monitor
Restart=on-failure
RestartSec=10

# Indispensable pour injecter les variables d'environnement graphiques
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority

[Install]
WantedBy=default.target
